# systemd/polymarket-strength-filtered.service
[Unit]
Description=Polymarket Strength-Filtered Monitor (>80% consensus only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rdent
Group=rdent
WorkingDirectory=/home/rdent/trading/polymarket

ExecStart=/usr/bin/python3 /home/rdent/trading/polymarket/scripts/strength-filtered-monitor.py

Restart=on-failure
RestartSec=10

Environment="PATH=/usr/local/bin:/usr/bin"
StandardOutput=append:/home/rdent/trading/polymarket/logs/strength-filtered-monitor.log
StandardError=append:/home/rdent/trading/polymarket/logs/strength-filtered-monitor.log

NoNewPrivileges=true
ProtectSystem=no
ProtectHome=no

[Install]
WantedBy=multi-user.target
